# data/systemd/mackes-ansible-pull.service
[Unit]
Description=Mackes Fleet — ansible-pull cycle
Documentation=man:ansible-pull(1) https://github.com/matthewmackes/MAP2-RELEASES
After=network-online.target headscale.service
Wants=network-online.target
ConditionPathExists=!/etc/mackes-shell/fleet.disabled

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 -m mackes.fleet --pull
ExecStartPost=/usr/bin/python3 -m mackes.fleet --prune
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
Nice=10

[Install]
WantedBy=multi-user.target
